# deploy/supervisord.conf
; Пример supervisord (Linux): положите в /etc/supervisord.d/ или подключите include.
; Пути и пользователя замените на свои. Рабочий каталог backend — где лежит app/.
[program:scalper-backend]
command=/opt/scalper/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
directory=/opt/scalper/backend
user=scalper
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stderr_logfile=/var/log/scalper/backend.err.log
stdout_logfile=/var/log/scalper/backend.out.log
environment=PATH="/opt/scalper/backend/.venv/bin"

[program:scalper-frontend]
command=/usr/bin/npm run preview -- --host 0.0.0.0 --port 5173
directory=/opt/scalper/frontend
user=scalper
autostart=true
autorestart=true
stderr_logfile=/var/log/scalper/frontend.err.log
stdout_logfile=/var/log/scalper/frontend.out.log

[program:scalper-watchdog]
command=/opt/scalper/backend/.venv/bin/python -m app.services.watchdog
directory=/opt/scalper/backend
user=scalper
autostart=false
autorestart=true
environment=WATCHDOG_ENABLED="true",PATH="/opt/scalper/backend/.venv/bin"
